# timers/mediator-csq-sync.timer
[Unit]
Description=Run CSQ (signal) sync of all meters

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
